# codex-lock-window-macos-fire.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.arnab.codex-lock-window.fire</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/arnab/Projects/scripts/codex-lock-window-macos</string>
    <string>run</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>19</integer>
    <key>Minute</key>
    <integer>30</integer>
    <key>Second</key>
    <integer>42</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/arnab/.local/state/codex-lock-window-macos/logs/fire.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/arnab/.local/state/codex-lock-window-macos/logs/fire.stderr.log</string>
</dict>
</plist>
